# files/etc/crontab
*/1 * * * * root env - `cat /root/env.txt` php /home/app/publish.php
