# files/etc/crontab
*/1 * * * * root php /home/app/publish.php